# com.my-remote.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.my-remote</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>/opt/homebrew/bin/cloudflared tunnel --config /Users/z5698891/Documents/my-remote/cloudflared-config.yml run my-remote &amp; /opt/homebrew/bin/node /Users/z5698891/Documents/my-remote/server.js &amp; wait</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/z5698891/Documents/my-remote</string>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/z5698891/Documents/my-remote/logs/server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/z5698891/Documents/my-remote/logs/server.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
